# data/cron/semaphore.cron
10 * * * * docker exec semaphore-core python3 src/semaphoreRunner.py -d test_dspec.json >> log.txt
20 */3 * * * docker exec semaphore-core python3 src/semaphoreRunner.py -d MLP-OP.json >> log.txt
