# data/cron/semaphore.cron
* * * * * docker exec semaphore-core python3 src/semaphoreRunner.py -d test_dspec.json >> log.txt
